# deploy/relay-server.service
[Unit]
Description=Relay WebSocket Server
After=network.target mysql.service
Wants=mysql.service

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/opt/relay
ExecStart=/opt/relay/relay-server
Restart=always
RestartSec=5
EnvironmentFile=/etc/relay/relay.env
StandardOutput=journal
StandardError=journal
SyslogIdentifier=relay-server

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/relay

[Install]
WantedBy=multi-user.target
